# roles/ansible-role-backup/templates/backup.service.j2
[Unit]
Description=Homelab Backup Job
After=network-online.target autofs.service
Wants=network-online.target autofs.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup.sh